# ops/rdq-research.service
# RD-Agent(Q) research control plane — systemd user unit (US-018).
#
# Runs the rdagent server_ui Flask control API (POST /trace, GET /receive,
# POST /user_interaction/submit, POST /control) on 127.0.0.1:19899 under the
# rdq-research OneCLI identity, so the research loops it spawns get
# Anthropic/Voyage/FMP auth injected by the proxy (http://127.0.0.1:10254).
#
# NOT exposed via tailscale serve — per the PLAN.md port table, server_ui has
# known flask-cors advisories and stays dark: the orchestrator talks to it
# over localhost only. The human-facing trace viewer (`rdagent ui`, :19900)
# is the only research UI ever exposed, and then tailnet-only (US-042);
# nothing on this box is ever exposed with `tailscale funnel`.
#
# Upstream `rdagent server_ui` binds 0.0.0.0; research/server_ui.py replicates
# its main() bound to loopback (the pinned upstream tree stays unmodified).
#
# Install:  ops/install_services.sh
# Start:    systemctl --user enable --now rdq-research.service
# Logs:     journalctl --user -u rdq-research.service -f

[Unit]
Description=RD-Agent(Q) research control plane (rdagent server_ui on 127.0.0.1:19899 via rdq-research OneCLI identity)
After=network-online.target
Wants=network-online.target
# Never park the unit in a permanent 'failed' state during longer gateway
# outages; RestartSec below paces the retries.
StartLimitIntervalSec=0

[Service]
Type=simple
WorkingDirectory=%h/rd-agent-q
# onecli run injects HTTPS_PROXY process-wide, and onecli preserves a pre-set
# NO_PROXY. Loopback control traffic (orchestrator -> :19899, spawned loops ->
# local mlflow/qlib services) and tailnet peers (100.64.0.0/10 CGNAT range)
# must not transit the OneCLI proxy.
Environment="NO_PROXY=127.0.0.1,localhost,100.64.0.0/10" "no_proxy=127.0.0.1,localhost,100.64.0.0/10"
# server_ui state lives outside the repo; rdagent.log.ui.conf reads these at
# import time (UI_ env prefix). research/server_ui.py creates both dirs.
Environment="UI_TRACE_FOLDER=%h/rdq-runs/server_ui/traces" "UI_STATIC_PATH=%h/rdq-runs/server_ui/static"

# --- US-market run environment (US-020) --------------------------------------
# fin_quant runs spawned via POST /upload inherit this environment; without it
# they silently backtest with rdagent's CN defaults (csi300, 2008-2020 dates).
# Values mirror ops/run_us_quant.sh wire_env — tests/test_services.py asserts
# the date defaults stay in sync with that script.
# Dates go under all three prefixes: the quant scenario prompt reads
# QLIB_QUANT_*, but the factor/model runners re-instantiate their own settings
# and render workspace conf jinja vars from QLIB_FACTOR_* / QLIB_MODEL_*.
Environment="QLIB_QUANT_TRAIN_START=2016-01-01" "QLIB_QUANT_TRAIN_END=2022-12-31" "QLIB_QUANT_VALID_START=2023-01-01" "QLIB_QUANT_VALID_END=2024-12-31" "QLIB_QUANT_TEST_START=2025-01-01" "QLIB_QUANT_TEST_END=2026-07-10"
Environment="QLIB_FACTOR_TRAIN_START=2016-01-01" "QLIB_FACTOR_TRAIN_END=2022-12-31" "QLIB_FACTOR_VALID_START=2023-01-01" "QLIB_FACTOR_VALID_END=2024-12-31" "QLIB_FACTOR_TEST_START=2025-01-01" "QLIB_FACTOR_TEST_END=2026-07-10"
Environment="QLIB_MODEL_TRAIN_START=2016-01-01" "QLIB_MODEL_TRAIN_END=2022-12-31" "QLIB_MODEL_VALID_START=2023-01-01" "QLIB_MODEL_VALID_END=2024-12-31" "QLIB_MODEL_TEST_START=2025-01-01" "QLIB_MODEL_TEST_END=2026-07-10"
# Factor coder source data (US-015 output) + US prompt overrides (US-016).
Environment="FACTOR_CoSTEER_DATA_FOLDER=%h/rdq-data/factor_source/us_liquid/data_folder"
Environment="FACTOR_CoSTEER_DATA_FOLDER_DEBUG=%h/rdq-data/factor_source/us_liquid/data_folder_debug"
Environment="APP_TPL=%h/rd-agent-q/research/app_tpl"
# US workspace templates re-injected via env-configurable class paths
# (APP_TPL cannot redirect the hardcoded template folders; research/us_quant.py).
Environment="QLIB_QUANT_FACTOR_HYPOTHESIS2EXPERIMENT=research.us_quant.USQlibFactorHypothesis2Experiment"
Environment="QLIB_QUANT_MODEL_HYPOTHESIS2EXPERIMENT=research.us_quant.USQlibModelHypothesis2Experiment"
# Execution environments (US-043): this box has docker but no conda, so qlib
# backtests/model training must opt out of rdagent's conda default (QTDockerEnv,
# image local_qlib:latest), and generated factor code runs with the repo venv's
# interpreter. Feature validation + the factor prompt env are runtime-patched
# to match (research/us_validation.py).
Environment="MODEL_CoSTEER_ENV_TYPE=docker"
Environment="FACTOR_CoSTEER_PYTHON_BIN=%h/rd-agent-q/.venv/bin/python"
# Never rebuild local_qlib:latest inside a run: upstream prepare() would run
# docker-py's LEGACY builder on every call (no BuildKit cache — a silent
# hour-long rebuild). Ops owns the image; rebuild after a pin bump with:
#   docker build -t local_qlib:latest <site-packages>/rdagent/scenarios/qlib/docker
Environment="QLIB_DOCKER_BUILD_FROM_DOCKERFILE=false"
# MLflow >= 3.6 (inside local_qlib:latest) raises MlflowException for the
# ./mlruns filesystem tracking backend unless this opt-out is set; qlib's
# recorder is hardwired to the file store, so every backtest dies in
# create_exp without it. Merged into the qlib container env via
# QlibDockerConf.env_dict (pydantic parses the JSON).
Environment="QLIB_DOCKER_ENV_DICT={\"MLFLOW_ALLOW_FILE_STORE\":\"true\"}"
# LiteLLM rejects temperature!=1 for the Claude models (rdagent defaults to
# chat_temperature=0.5) — drop unsupported params instead of erroring
# (research/CLAUDE.md; the first fin_quant chat call dies without this).
Environment="LITELLM_DROP_PARAMS=true"
# Voyage-backed CoSTEER RAG (research/us_validation.py shim #5). Off by
# default because Voyage's free tier rate-limited runs to death; a payment
# method was added 2026-07-10, so real embeddings + cross-loop knowledge
# accumulation are back on. Unset (or remove) to fall back to the local
# deterministic-embedding shim if Voyage misbehaves again.
Environment="RDQ_ENABLE_RAG=1"
# Keep spawned-run workspaces out of the repo, next to the trace state.
Environment="WORKSPACE_PATH=%h/rdq-runs/server_ui/workspace"
# LLM backend (CHAT_MODEL/EMBEDDING_MODEL + placeholder keys the proxy
# overrides). Deliberately NOT optional: a missing file must fail the unit
# loudly instead of letting spawned runs fall back to rdagent's OpenAI
# defaults. Copy research/.env.example if this stops the service.
EnvironmentFile=%h/rd-agent-q/research/.env

ExecStart=%h/.local/bin/onecli run --agent rdq-research -- %h/rd-agent-q/.venv/bin/python -m research.server_ui
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
